# overlay/etc/systemd/user/update-tools.service
[Unit]
Description=Update Rust, Node, Python, and Rust CLIs (per-user)

[Service]
Type=oneshot
ExecStart=/usr/local/libexec/update-tools.sh
