# overlay/etc/systemd/system/user/blueman-applet.service
[Unit]
Description=Blueman Bluetooth applet
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/blueman-applet

[Install]
WantedBy=graphical-session.target
