# overlay/etc/systemd/system/user/chelipeds-devsetup.service
[Unit]
Description=Per-user dev tool bootstrap (rustup/uv/volta + Rust CLIs)

[Service]
Type=oneshot
ExecStart=/usr/local/libexec/chelipeds-devsetup.sh
RemainAfterExit=yes
